# tod.cron
# Time of Day taskwarrior extension cron file
# requires 'ln -s ./tod.cron /etc/cron.d/tod.cron
# and permissions set
#
# This file should be generated by ToD.sh, somehow..

# TODO get redirect messages to /dev/null right !
* 8 * * * bash 'ToD morn >$%something ?dev/nul'
* 8 * * * bash 'ToD day >$%something ?dev/nul'
* 13 * * * bash 'ToD aft >$%something ?dev/nul'
* 18 * * * bash 'ToD eve >$%something ?dev/nul'
* 21 * * * bash 'ToD night >$%something ?dev/nul'

